# nightly_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  VulnScanner 改善ループ — Windows タスクスケジューラ設定ファイル（参考用）

  ※ 現在は手動実行のみ。スケジューラ登録は不要。
  ※ 手動実行: .\run_improvement.ps1 [-MaxHours 4] [-StopAt "05:00"]

  スケジューラに登録したい場合のインポート手順:
    1. 環境変数を設定しておく:
         FUGU_API_KEY=<your-key>
         FUGU_API_BASE=https://api.fugu.sakana.ai/v1   (省略可)
         FUGU_MODEL=fugu-chat                          (省略可)
    2. 管理者権限の PowerShell で実行:
         Register-ScheduledTask -Xml (Get-Content nightly_task.xml -Raw) -TaskName "VulnScannerImprovement" -Force
    3. 確認:
         Get-ScheduledTask -TaskName "VulnScannerImprovement"

  削除:
    Unregister-ScheduledTask -TaskName "VulnScannerImprovement" -Confirm:$false

  注意: <Command> の python パスは環境に合わせて変更してください。
        フルパスを使う場合の例:
          C:\Users\dorom\AppData\Local\Programs\Python\Python314\python.exe
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>VulnScanner 改善ループ — 改善提案を生成・評価する。</Description>
    <Author>VulnScanner</Author>
    <URI>\VulnScannerImprovement</URI>
  </RegistrationInfo>

  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-06-29T02:00:00</StartBoundary>
      <Enabled>false</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>false</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT5H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <Command>python</Command>
      <Arguments>improvement_loop.py --max-hours 4</Arguments>
      <WorkingDirectory>C:\VulnScanner</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
